# scripts/ecfr-sync.timer
[Unit]
Description=Monthly eCFR regulatory sync for TTB label verifier

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
